# portmap/com.apple.portmap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.portmap</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/sandbox-exec</string>
		<string>-f</string>
		<string>/usr/share/sandbox/portmap.sb</string>
		<string>/usr/sbin/portmap</string>
	</array>
	<key>Sockets</key>
	<dict>
		<key>OnDemandTickler</key>
		<dict>
			<key>SockPathName</key>
			<string>/var/run/portmap.socket</string>
			<key>SockPathMode</key>
			<integer>511</integer>
		</dict>
	</dict>
	<key>ServiceIPC</key>
	<true/>
	<key>HopefullyExitsLast</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
</dict>
</plist>
