# Security/OSX/sec/ipc/com.apple.securityd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>DEBUGSCOPE</key>
		<string>-policy-node,policy-set,policy,alloc,trust,bind,profile,trace,dbconn,OTR,serverxpc,sqlite3,error_thee_well, ringSigning</string>
		<key>WAIT4DEBUGGER</key>
		<string>NO</string>
	</dict>
	<key>GroupName</key>
	<string>_securityd</string>
	<key>Label</key>
	<string>com.apple.securityd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.securityd</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/securityd</string>
	</array>
	<key>Umask</key>
	<integer>54</integer>
	<key>UserName</key>
	<string>_securityd</string>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>MinimalBootProfile</key>
	<true/>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>kSOSCCHoldLockForInitialSync</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.security.secureobjectsync.holdlock</string>
			</dict>
		</dict>
	</dict>
</dict>
</plist>
