# KerberosLibraries/kadmind.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>edu.mit.Kerberos.kadmind</string>
	<key>Program</key>
	<string>/usr/sbin/kadmind</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/sbin/kadmind</string>
		<string>-nofork</string>
	</array>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/Library/Preferences/edu.mit.Kerberos.kadmind.launchd</key>
			<true/>
		</dict>
	</dict>
	<key>EnableTransactions</key>
	<true/>
</dict>
</plist>
