# autofs/autofsd/com.apple.autofsd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.autofsd</string>
	<key>KeepAlive</key>
	<true/>
	<key>Program</key>
	<string>/usr/libexec/autofsd</string>
	<key>ProgramArguments</key>
	<array>
		<string>autofsd</string>
	</array>
	<key>EnableTransactions</key>
	<true/>
</dict>
</plist>
