# gssd/com.apple.gssd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>com.apple.gssd</string>
        <key>MachServices</key>
        <dict>
                <key>com.apple.gssd</key>
                <dict>
                        <key>TaskSpecialPort</key>
                        <integer>8</integer>
                </dict>
        </dict>
        <key>ProgramArguments</key>
        <array>
                <string>/usr/sbin/gssd</string>
        </array>
        <key>EnvironmentVariables</key>
        <dict>
                <key>CFFIXED_USER_HOME</key>
                <string>/tmp</string>
        </dict>
        <key>HopefullyExitsLast</key>
        <true/>
        <key>EnableTransactions</key>
        <true/>        
</dict>
</plist>
