# mailman/Mailman.LaunchDaemons/org.list.mailmanctl.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Debug</key>
	<true/>
	<key>Disabled</key>
	<true/>
	<key>Label</key>
	<string>org.list.mailmanctl</string>
	<key>OnDemand</key>
	<false/>
	<key>Program</key>
	<string>/usr/share/mailman/bin/mailmanctl</string>
	<key>ProgramArguments</key>
	<array>
		<string>mailmanctl</string>
		<string>-s</string>
		<string>startf</string>
	</array>
	<key>ServiceIPC</key>
	<false/>
</dict>
</plist>
