# system_cmds/dynamic_pager.tproj/com.apple.dynamic_pager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
                <key>Label</key>
                <string>com.apple.dynamic_pager</string>
                <key>ProgramArguments</key>
                <array>
                     <string>/sbin/dynamic_pager</string>
                     <string>-F</string>
                     <string>/private/var/vm/swapfile</string>
                </array>
                <key>OnDemand</key>
                <false/>
		<key>HopefullyExitsLast</key>
		<true/>
		<key>EnableTransactions</key>
		<true/>
</dict>
</plist>
