# Security/CloudKeychainProxy/com.apple.security.cloudkeychainproxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>LimitLoadToSessionType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/Security.framework/Versions/A/Resources/CloudKeychainProxy.bundle/Contents/MacOS/CloudKeychainProxy</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.security.cloudkeychainproxy3</key>
		<true/>
	</dict>
	<key>RunAtLoad</key>
	<false/>
	<key>KeepAlive</key>
	<false/>
	<key>Program</key>
	<string>/System/Library/Frameworks/Security.framework/Versions/A/Resources/CloudKeychainProxy.bundle/Contents/MacOS/CloudKeychainProxy</string>
	<key>Label</key>
	<string>com.apple.security.cloudkeychainproxy3</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>DEBUGSCOPE</key>
		<string>all</string>
		<key>WAIT4DEBUGGER</key>
		<string>NO</string>
	</dict>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.keystore.lockstatus</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.keystore.lockstatus</string>
			</dict>
			<key>com.apple.security.cloudkeychainproxy.kvstorechange3</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.security.cloudkeychainproxy.kvstorechange3</string>
			</dict>
			<key>com.apple.security.cloudkeychain.forceupdate</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.security.cloudkeychain.forceupdate</string>
			</dict>
		</dict>
	</dict>
</dict>
</plist>
