# PowerManagement/swd/com.apple.powerd.swd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>ThrottleInterval</key>
    <integer>30</integer>
	<key>Label</key>
	<string>com.apple.powerd.swd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/powerd.bundle/swd</string>
	</array>
    <key>WatchPaths</key>
    <array>
        <string>/var/tmp/SleepWakeLog.dump</string>
        <string>/var/tmp/PMSleepWakeLog.dump</string>
    </array>
</dict>
</plist>
