# smb/launchd/com.apple.smbfs_load_kext.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.smbfs_load_kext</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.smbfs_load_kext</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Extensions/smbfs.kext/Contents/Resources/load_smbfs</string>
	</array>
</dict>
</plist>
