# distcc/distccd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>GroupName</key>
	<string>nobody</string>
	<key>Label</key>
	<string>com.apple.distccd</string>
	<key>OnDemand</key>
	<false/>
	<key>Program</key>
	<string>/usr/bin/distccd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/distccd</string>
		<string>--no-detach</string>
		<string>--daemon</string>
		<string>--priority</string>
		<string>10</string>
		<string>--min-disk-free</string>
		<string>512</string>
	</array>
	<key>SHAuthorizationRight</key>
	<string>com.apple.Xcode.distcc.admin</string>
	<key>RunAtLoad</key>
	<true/>
	<key>ServiceDescription</key>
	<string>Distributed compilation daemon (distccd)</string>
	<key>UserName</key>
	<string>nobody</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>USE_XCODE_SELECT_PATH</key>
		<string>1</string>
	</dict>
</dict>
</plist>
