# NFS/rpc.lockd/com.apple.lockd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>com.apple.lockd</string>
        <key>MachServices</key>
        <dict>
                <key>com.apple.lockd</key>
                <dict>
                        <key>HostSpecialPort</key>
                        <integer>12</integer>
                </dict>
        </dict>
        <key>ProgramArguments</key>
        <array>
                <string>/usr/sbin/rpc.lockd</string>
        </array>
        <key>HopefullyExitsLast</key>
        <true/>
        <key>KeepAlive</key>
        <dict>
                <key>PathState</key>
		<dict>
			<key>/var/run/lockd.pid</key>
			<true/>
		</dict>
        </dict>
</dict>
</plist>
