# crontabs/private/etc/crontab
# /etc/crontab
SHELL=/bin/sh
PATH=/etc:/bin:/sbin:/usr/bin:/usr/sbin
HOME=/var/log
#
#minute	hour	mday	month	wday	who	command
#
#*/5	*	*	*	*	root	/usr/libexec/atrun
#
# Run daily/weekly/monthly jobs.
15	3	*	*	*	root	periodic daily
30	4	*	*	6	root	periodic weekly
30	5	1	*	*	root	periodic monthly
