# Security/sec/ipc/com.apple.secd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.secd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.secd</key>
		<true/>
	</dict>
	<key>OnDemand</key>
	<true/>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/secd</string>
	</array>
	<key>ServiceIPC</key>
	<true/>
	<key>Umask</key>
	<integer>54</integer>
</dict>
</plist>
